# backend-swift/packaging/wiki.qaq.unfaird.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>wiki.qaq.unfaird</string>
    <key>ProgramArguments</key>
    <array>
        <string>@INSTALL_PREFIX@/usr/local/lib/unfaird/UnfairDaemon</string>
        <string>serve</string>
        <string>--hostname</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8080</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>POSIXSpawnType</key>
    <string>Interactive</string>
    <key>JetsamProperties</key>
    <dict>
        <key>JetsamMemoryLimit</key>
        <integer>512</integer>
        <key>JetsamPriority</key>
        <integer>-49</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/var/tmp/unfaird</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DYLD_LIBRARY_PATH</key>
        <string>@DYLD_LIBRARY_PATH@</string>
        <key>PATH</key>
        <string>@LAUNCHD_PATH@</string>
        <key>PUBLIC_DIR</key>
        <string>@INSTALL_PREFIX@/usr/share/assppweb/public</string>
        <key>DATA_DIR</key>
        <string>/var/mobile/AssppWebData</string>
        <key>UNSAFE_DANGEROUSLY_DISABLE_HTTPS_REDIRECT</key>
        <string>true</string>
        <key>DOWNLOAD_THREADS</key>
        <string>8</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/var/tmp/unfaird/logs/unfaird.log</string>
    <key>StandardErrorPath</key>
    <string>/var/tmp/unfaird/logs/unfaird.err</string>
</dict>
</plist>
